# systemd/user/arrpc.service
[Unit]
Description=Discord Rich Presence Client

[Service]
Restart=always
ExecStart=/usr/bin/arrpc

[Install]
WantedBy=default.target