# systemd/user/qpwgraph.service
[Unit]
Description=Start QPWGraph

[Service]
Restart=always
ExecStart=/usr/bin/qpwgraph -a -m

[Install]
WantedBy=default.target