# systemd/user/restore-hypr-session.service
[Unit]
Description=Restore Hyprland session after login
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=/home/dei/.config/hypr/session/restore.sh

[Install]
WantedBy=graphical-session.target
