# systemd/user/waybar.service
[Unit]
Description=Waybar Check
[Service]
ExecStart=/home/dei/.config/waybar/scripts/check_waybar.sh