# systemd/user/waybar.timer
[Unit]
Description=Check Waybar
[Timer]
OnBootSec=10
OnUnitActiveSec=10
AccuracySec=1ms
[Install]
WantedBy=timers.target