# systemd/user/waybar.timer
[Unit]
Description=Check Waybar
[Timer]
OnBootSec=5
OnUnitActiveSec=1
AccuracySec=1ms
[Install]
WantedBy=timers.target